# 09_deliverables/board_deploy/systemd/smartbag-cloud-uploader.service
[Unit]
Description=Optional SS928 SmartBag CloudBase telemetry uploader
After=network-online.target smartbag-alert.service
Wants=network-online.target

[Service]
Type=simple
User=root
EnvironmentFile=-/etc/smartbag/cloud-uploader.env
WorkingDirectory=/root/smartbag/cloud-uploader
ExecStart=/usr/bin/python3 /root/smartbag/cloud-uploader/cloud_uploader.py --config /etc/smartbag/cloud-uploader.json
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
