# 09_deliverables/board_deploy/systemd/smartbag-video.service
[Unit]
Description=SS928 SmartBag dual-camera LAN video gateway
After=network.target
Conflicts=smartbag-controller.service smartbag-alternating-vision.service

[Service]
Type=simple
WorkingDirectory=/root/smartbag/vision
EnvironmentFile=-/etc/smartbag/smartbag.env
ExecStart=/root/smartbag/venv/bin/python /root/smartbag/vision/dual_camera_gateway.py --config /etc/smartbag/config.json
Restart=on-failure
RestartSec=2
KillMode=control-group
TimeoutStopSec=10
StandardInput=null
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
