# 09_deliverables/board_deploy/systemd/smartbag-controller.service
[Unit]
Description=SS928 SmartBag autonomous controller and supervised runtime
After=local-fs.target smartbag-safe-off.service bluetooth.service
Wants=smartbag-safe-off.service bluetooth.service
Conflicts=smartbag-alert.service smartbag-alternating-vision.service smartbag-video.service

[Service]
Type=simple
User=root
RuntimeDirectory=smartbag
RuntimeDirectoryMode=0755
WorkingDirectory=/root/smartbag/controller
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/smartbag/smartbag.env
ExecStartPre=/root/smartbag/apply-smartbag-pinmux.sh
ExecStartPre=/root/smartbag/venv/bin/python /root/smartbag/board-deploy/wait_for_hardware.py --profile controller --config /etc/smartbag/config.json --timeout-s 30
ExecStart=/root/smartbag/venv/bin/python /root/smartbag/controller/smartbag_alert_controller.py --config /etc/smartbag/config.json --detector-cwd /root/smartbag/vision --gnss-command "/root/smartbag/venv/bin/python /root/smartbag/gnss/dx_gp21_tracker.py --config /root/smartbag/gnss/config.ss928_uart4.json --command-stdin --no-ble" --imu-command "/root/smartbag/venv/bin/python /root/smartbag/imu/bmi270_backpack.py --config /root/smartbag/imu/config.example.json --hardware-profile /etc/smartbag/hardware.json --command-stdin --no-ble"
ExecStopPost=-/root/smartbag/venv/bin/python /root/smartbag/board-deploy/safe_off.py --hardware /etc/smartbag/hardware.json
Restart=always
RestartSec=2
KillMode=control-group
TimeoutStopSec=12
StandardInput=null
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=smartbag.target
