# 09_deliverables/board_deploy/systemd/smartbag-safe-off.service
[Unit]
Description=SS928 SmartBag startup and shutdown actuator safe-off
After=local-fs.target
Before=smartbag-controller.service shutdown.target reboot.target halt.target

[Service]
Type=oneshot
WorkingDirectory=/root/smartbag/board-deploy
EnvironmentFile=-/etc/smartbag/smartbag.env
ExecStart=/root/smartbag/venv/bin/python /root/smartbag/board-deploy/safe_off.py --hardware /etc/smartbag/hardware.json
ExecStop=/root/smartbag/venv/bin/python /root/smartbag/board-deploy/safe_off.py --hardware /etc/smartbag/hardware.json
RemainAfterExit=yes
TimeoutStartSec=10
TimeoutStopSec=10
StandardInput=null
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=smartbag.target
