# 09_deliverables/board_deploy/systemd/smartbag-alternating-cleanup.timer
[Unit]
Description=Periodic SS928 SmartBag alternating-camera session cleanup

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
