# 09_deliverables/board_deploy/systemd/smartbag.target
[Unit]
Description=SS928 SmartBag Rev2 autonomous power-only runtime
Requires=smartbag-safe-off.service smartbag-controller.service
Wants=smartbag-boot-selftest.service
After=smartbag-safe-off.service smartbag-controller.service

[Install]
WantedBy=multi-user.target
